# packages/xbox-switch-bridge/xbox_switch_bridge-0.3.3.tar.gz/xbox_switch_bridge-0.3.3/data/systemd/xbox-switch-bridge.service
[Unit]
Description=Xbox to Nintendo Switch Controller Bridge
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/xbox-switch-bridge --run
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/xbox-switch-bridge.log
StandardError=append:/var/log/xbox-switch-bridge.log

[Install]
WantedBy=multi-user.target