# packages/utilmy/utilmy-0.1.17367756-py3-none-any.whl/utilmy/asearch/rag/engines/diffbot_llm_kg_jscript/supervisord.conf
[supervisord]
nodaemon=true

[program:diffbot-llm-inference]
command=/bin/bash -c "cd /code && poetry run sh ./start_server.sh 2>&1 | tee -a /var/log/diffbot.log | cat";
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:vllm]
command=/bin/bash -c "python3 -m vllm.entrypoints.openai.api_server ${VLLM_OPTIONS} 2>&1 | tee -a /var/log/vllm.log | cat"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0